# virttmp/rc.local
#!/bin/bash
# THIS FILE IS ADDED FOR COMPATIBILITY PURPOSES
#
# It is highly advisable to create own systemd services or udev rules
# to run scripts during boot instead of using this file.
#
# In constrast to previous versions due to parallel execution during boot 
# this script will NOT be run after all other services.
#  
# Please note that you must run 'chmod +x /etc/rc.d/rc.local' to ensure
# that this script will be executed during boot.

touch /var/lock/subsys/local
sh /bin/setipbymac.sh centos7-190 10.20.102.1 10.20.102.190 255.255.255.0 52:54:00:EA:AB:77 172.16.2.190 255.255.255.0 52:54:00:28:16:B7
